# examples/metapi.service
[Unit]
Description=metapi-rs lightweight LLM router
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=metapi
Group=metapi
WorkingDirectory=/opt/metapi
EnvironmentFile=/etc/metapi.env
ExecStart=/opt/metapi/metapi-rs
Restart=always
RestartSec=3
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
